# airootfs/etc/systemd/system/suyuos-build.service
[Unit]
Description=suyuOS Emulator Build Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/build-emulators
StandardOutput=journal
StandardError=journal
TimeoutStartSec=3600
User=root

[Install]
WantedBy=multi-user.target